# deploy/systemd/transcribe-frontend-dev.service
[Unit]
Description=Transcribe Dev Frontend + /api Proxy
After=transcribe-api-dev.service
Requires=transcribe-api-dev.service

[Service]
EnvironmentFile=-/home/gunnar/.config/transcribe/dev.env
Type=simple
WorkingDirectory=/home/gunnar/projects/transcribe-dev
ExecStart=/usr/bin/python3 /home/gunnar/projects/transcribe-dev/deploy/dev_frontend_proxy.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
