# deploy/systemd/transcribe-worker-dev@.service
# User unit (systemctl --user). Template so you can run multiple workers: transcribe-worker-dev@1, @2, ...
[Unit]
Description=Transcribe Dev Worker (%i)

[Service]
Type=simple
WorkingDirectory=/home/gunnar/projects/transcribe-dev/worker
ExecStartPre=/usr/bin/mkdir -p /home/gunnar/projects/transcribe-dev/data/demo_jobs/inbox /home/gunnar/projects/transcribe-dev/data/demo_jobs/running /home/gunnar/projects/transcribe-dev/data/demo_jobs/done /home/gunnar/projects/transcribe-dev/data/demo_jobs/error
ExecStart=/srv/transcribe/worker/.venv/bin/python /home/gunnar/projects/transcribe-dev/worker/worker_daemon.py
Restart=always
RestartSec=2
EnvironmentFile=-%h/.config/transcribe/dev.env
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=default.target
