# deploy/systemd/transcribe-frontend-dev.service
[Unit]
Description=Transcribe Dev Frontend + /api Proxy
After=transcribe-api-dev.service
Requires=transcribe-api-dev.service

[Service]
Type=simple
WorkingDirectory=/home/gunnar/projects/transcribe-dev
ExecStart=/usr/bin/python3 /home/gunnar/projects/transcribe-dev/deploy/dev_frontend_proxy.py
Environment=TRANSCRIBE_FRONTEND_DIR=/home/gunnar/projects/transcribe-dev/static
Environment=TRANSCRIBE_API_BASE=http://127.0.0.1:8001
Environment=TRANSCRIBE_FRONTEND_HOST=127.0.0.1
Environment=TRANSCRIBE_FRONTEND_PORT=8010
Environment=TRANSCRIBE_FRONTEND_WS_BUFFER_BYTES=65536
Environment=TRANSCRIBE_FRONTEND_WS_IDLE_TIMEOUT_S=1800
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
