# deploy/systemd/transcribe-whisperlive-dev.service
[Unit]
Description=WhisperLive Sidecar (Dev, Docker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker rm -f whisperlive-sidecar
ExecStart=/usr/bin/docker run --rm --name whisperlive-sidecar --gpus all -p 127.0.0.1:9090:9090 ghcr.io/collabora/whisperlive-gpu:latest
ExecStop=/usr/bin/docker stop whisperlive-sidecar
Restart=always
RestartSec=3
TimeoutStopSec=20

[Install]
WantedBy=default.target
