# deploy/systemd/transcribe-api-dev.service
[Unit]
Description=Transcribe Dev API (Uvicorn)

[Service]
Type=simple
WorkingDirectory=/home/gunnar/projects/transcribe-dev/portal-api
ExecStart=/srv/transcribe/portal-api/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
